# services.d/privoxy/run
#!/bin/bash -ex

exec /usr/sbin/privoxy --no-daemon /etc/privoxy/config 2>&1
