# services.d/tor/finish
#!/bin/bash -ex
#/usr/bin/s6-svscanctl -t /etc/services.d/tor
/bin/true
