# services.d/tor/run
#!/bin/bash -ex

exec /usr/bin/tor -f /etc/tor/torrc --hush 2>&1
